# templates/default/netcalc-upstart.conf.erb
start on runlevel [2345]
stop on runlevel [016]

exec su serf -c 'python /home/tornado/netcalc/app/netcalc_listener.py'
